# .docker/php-fpm/crontab
* * * * * root php /var/www/bin/console notify:mrs-to-teams >> /var/log/cron.log 2>&1
